# provisioning/ansible/roles/app_systemd/templates/app.target.j2
[Unit]
StopWhenUnneeded=true
Wants=app-{{app_user}}-web.target

[Install]
WantedBy=multi-user.target
