# provisioning/ansible/roles/app_unicorn/templates/app-web.target.j2
[Unit]
StopWhenUnneeded=true
Wants=app-{{app_user}}-web-1.service
